# install_dir/service/mon_super_archiver.service
[Unit]
Description=Service de lancement automatique du script d'archivage

[Service]
ExecStart=/bin/bash /usr/bin/archiver/main.sh
User=$USER