# install_dir/service/mon_super_archiver.service
[Unit]
Description=Service de lancement du script d'archivage

[Service]
Type=simple
ExecStart=/bin/bash /bin/mon_super_archiver/main.sh

